# gentoo/.config/systemd/user/redshift.service.d/custom.conf
[Service]
Environment=DISPLAY=:0
RestartSec=20
