# gentoo/etc/systemd/system/powertop.service
[Unit]
Description=Powertop tunings

[Service]
Type=oneshot
ExecStart=/usr/local/bin/manual-powertop

[Install]
WantedBy=multi-user.target
